# supervisor.conf
[program:main]
command=/bin/bash /workspace/vllm/check_runpod.sh
directory=/workspace/vllm/
autostart=true
autorestart=true
startsecs = 3
stderr_logfile=/vllm.err.log
stdout_logfile=/vllm.out.log
